# deploy/systemd/qevor-api.service
[Unit]
Description=Qevor API (Express)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=qevor
Group=qevor
WorkingDirectory=/opt/qevor/server
EnvironmentFile=/etc/qevor/qevor-api.env
ExecStart=/usr/bin/node /opt/qevor/server/dist/api/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/log/qevor
ProtectHome=true

[Install]
WantedBy=multi-user.target
